# debian/wmbeacons.cron.weekly
#!/bin/sh

cd /srv/wmbeacons
carton exec bin/wpextract de dewiki.templates.gz
carton exec bin/wpbeacons de
